# src/g00x/systemd_templates/box.timer
[Unit]
Description=Sync Box mnt to fsx user workspace

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
